# claude-telegram-bot.service
[Unit]
Description=Claude Telegram Bot
After=network.target

[Service]
Type=simple
# Change these paths to match your installation
User=YOUR_USERNAME
WorkingDirectory=/path/to/claude-telegram-bot
EnvironmentFile=/path/to/claude-telegram-bot/.env
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/claude-telegram-bot/venv/bin/python /path/to/claude-telegram-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
